# roles/system_tty_login/templates/hyprland-session.service.j2
[Unit]
Description=Hyprland session

[Service]
Type=simple
ExecStart=%h/.local/bin/start-hyprland
Restart=no

[Install]
WantedBy=default.target
